# deploy/aegis-funder-sync.service
[Unit]
Description=AEGIS — sync funder guidelines from Google Drive
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=aegis
Group=aegis
WorkingDirectory=/opt/aegis
EnvironmentFile=/etc/aegis/aegis.env
ExecStart=/opt/aegis/.venv/bin/python scripts/sync_funders_from_folder.py --apply
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aegis-funder-sync

# Same hardening posture as aegis-web / aegis-worker.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/aegis
